# docker/rootfs/etc/s6-overlay/s6-rc.d/reverb/run
#!/command/with-contenv sh
exec s6-setuidgid sendtrap php /app/artisan reverb:start --host=127.0.0.1 --port=9090
