# docker/rootfs/etc/s6-overlay/s6-rc.d/bootstrap/up
#!/command/execlineb -P
/usr/local/bin/sendtrap-bootstrap
