# docker/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv sh
# The FPM master opens the inherited container log descriptors, then the pool
# configuration drops every request worker to `sendtrap`.
exec php-fpm -F
